# haven/systemd/lyra-haven.service
[Unit]
Description=Lyra Haven Bot - Haven chat presence and conversation handling
Documentation=file:///mnt/c/Users/Jeff/Claude_Projects/Awareness/haven/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

WorkingDirectory=/mnt/c/Users/Jeff/Claude_Projects/Awareness

Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONDONTWRITEBYTECODE=1"
Environment="ENTITY_NAME=lyra"
Environment="ENTITY_TOKEN_FILE=/mnt/c/Users/Jeff/Claude_Projects/Awareness/entities/lyra/.entity_token"
Environment="ENTITY_PATH=/mnt/c/Users/Jeff/Claude_Projects/Awareness/entities/lyra"
Environment="HAVEN_URL=http://localhost:8205"
Environment="PPS_HTTP_URL=http://localhost:8201"
Environment="CLAUDE_MODEL=sonnet"
Environment="DEBOUNCE_JITTER_SECONDS=0"
Environment="ALWAYS_RESPOND=1"
Environment="HAVEN_GATE_ENABLED=1"
Environment="PROJECT_DIR=/mnt/c/Users/Jeff/Claude_Projects/Awareness"

ExecStart=/mnt/c/Users/Jeff/Claude_Projects/Awareness/.venv/bin/python3 -m haven.bot

Restart=always
RestartSec=15
StartLimitInterval=300
StartLimitBurst=5

KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

MemoryMax=1024M
MemorySwapMax=512M
CPUQuota=50%

StandardOutput=journal
StandardError=journal
SyslogIdentifier=lyra-haven

[Install]
WantedBy=default.target
