# daemon/systemd/lyra-backup.service
[Unit]
Description=Lyra PPS Backup - Daily backup of all entity data and shared databases
After=docker.service

[Service]
Type=oneshot
WorkingDirectory=/mnt/c/Users/Jeff/Claude_Projects/Awareness
ExecStart=/mnt/c/Users/Jeff/Claude_Projects/Awareness/.venv/bin/python3 scripts/backup_pps.py --no-stop
StandardOutput=append:/tmp/lyra_backup.log
StandardError=append:/tmp/lyra_backup.log

[Install]
WantedBy=default.target
