# daemon/systemd/lyra-reflection.service
[Unit]
Description=Lyra Reflection Daemon - Autonomous reflection and maintenance
Documentation=file:///mnt/c/Users/Jeff/Claude_Projects/Awareness/daemon/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# Note: User/Group not needed for user services, and cause issues in WSL2

WorkingDirectory=/mnt/c/Users/Jeff/Claude_Projects/Awareness/daemon

EnvironmentFile=/mnt/c/Users/Jeff/Claude_Projects/Awareness/daemon/.env

Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONDONTWRITEBYTECODE=1"

# Use project-level venv (consolidated in Issue #111)
ExecStart=/mnt/c/Users/Jeff/Claude_Projects/Awareness/.venv/bin/python3 -u lyra_reflection.py

Restart=on-failure
RestartSec=30
StartLimitInterval=300
StartLimitBurst=3

KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=60

MemoryMax=512M
MemorySwapMax=512M
CPUQuota=25%

# Security settings disabled for WSL2 compatibility
# NoNewPrivileges=true
# PrivateTmp=true
# ProtectSystem=strict
# ProtectHome=no
# ReadWritePaths=/mnt/c/Users/Jeff/Claude_Projects/Awareness /home/jeff/.claude

StandardOutput=journal
StandardError=journal
SyslogIdentifier=lyra-reflection

[Install]
WantedBy=default.target
