# daemon/systemd/lyra-discord.service
[Unit]
Description=Lyra Discord Daemon - Discord presence and conversation handling
Documentation=file:///mnt/c/Users/Jeff/Claude_Projects/Awareness/daemon/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# Note: User/Group not needed for user services, and cause issues in WSL2

WorkingDirectory=/mnt/c/Users/Jeff/Claude_Projects/Awareness/daemon

EnvironmentFile=/mnt/c/Users/Jeff/Claude_Projects/Awareness/daemon/.env

Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONDONTWRITEBYTECODE=1"

# Use project-level venv (consolidated in Issue #111)
ExecStart=/mnt/c/Users/Jeff/Claude_Projects/Awareness/.venv/bin/python3 -u lyra_daemon.py

Restart=on-failure
RestartSec=10
StartLimitInterval=60
StartLimitBurst=5

KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

MemoryMax=1024M
MemorySwapMax=512M
CPUQuota=50%

# Security settings disabled for WSL2 compatibility
# NoNewPrivileges=true
# PrivateTmp=true
# ProtectSystem=strict
# ProtectHome=no
# ReadWritePaths=/mnt/c/Users/Jeff/Claude_Projects/Awareness/daemon/logs /home/jeff/.claude

StandardOutput=journal
StandardError=journal
SyslogIdentifier=lyra-discord

[Install]
WantedBy=default.target
